# elkscmd/rootfs_template/etc/inittab
id:3:initdefault:

si::sysinit:/etc/rc.d/rc.sys

#l0:0:wait:/etc/rc.d/rc 0
#l1:1:wait:/etc/rc.d/rc 1
#l2:2:wait:/etc/rc.d/rc 2
#l3:3:wait:/etc/rc.d/rc 3
#l6:6:wait:/etc/rc.d/rc 6

#ud::once:/sbin/update

1:2345:respawn:/bin/getty /dev/tty1
#2:2345:respawn:/bin/getty /dev/tty2
#3:2345:respawn:/bin/getty /dev/tty3
